# scripts/clipboard-monitor.service
# systemd user service for Memory Palace clipboard monitor (Linux)
#
# TEMPLATE — edit paths below before installing:
#   Replace /home/YOUR_USER/Projects/memory-palace with your actual project path.
#   Replace YOUR_HOSTNAME with a short identifier for this machine.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp clipboard-monitor.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now clipboard-monitor
#
# Logs: journalctl --user -u clipboard-monitor -f

[Unit]
Description=Memory Palace clipboard monitor
Documentation=https://github.com/kashfshah/memory-palace
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/python3 /home/YOUR_USER/Projects/memory-palace/scripts/clipboard-monitor.py
WorkingDirectory=/home/YOUR_USER/Projects/memory-palace
Restart=always
RestartSec=10

Environment=DISPLAY=:0
Environment=CLIPBOARD_MACHINE=YOUR_HOSTNAME
Environment=CLIPBOARD_JSONL_PATH=/home/YOUR_USER/Projects/memory-palace/data/clipboard.jsonl

StandardOutput=append:/home/YOUR_USER/.local/share/clipboard-monitor.log
StandardError=append:/home/YOUR_USER/.local/share/clipboard-monitor.log

[Install]
WantedBy=graphical-session.target
