# scripts/net.kashifshah.clipboard-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template — run scripts/install-clipboard-monitor.sh to install with correct paths. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>net.kashifshah.clipboard-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON3__</string>
    <string>__PROJECT_DIR__/scripts/clipboard-monitor.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CLIPBOARD_MACHINE</key>
    <string>__HOSTNAME__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/clipboard-monitor.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/clipboard-monitor.log</string>
</dict>
</plist>
